# files/collectd-supervisor.conf
[program:collectd]
autostart=true
autorestart=true
command=/usr/bin/confd -onetime -verbose -backend env -confdir /etc/confd && /usr/sbin/collectd -f