# files/collectd-supervisor.conf
[program:collectd]
autorestart=true
command=/opt/start-collectd.sh
